# recipes-core/dnf-automatic-restart/dnf-automatic-restart.service

[Unit]
Description=dnf automatic install updates

[Service]
Type=simple
Environment="ABRT_IGNORE_PYTHON=1"
ExecStart=/usr/bin/dnf-automatic-restart
