# recipes-config/dropbear/dropbear.socket
[Unit]
Conflicts=dropbear.service

[Socket]
ListenStream=41415
Accept=yes

[Install]
WantedBy=sockets.target
Also=dropbearkey.service
